# charybdis/rootfs/etc/services.d/ngircd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the ngircd service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# shellcheck shell=bash

bashio::log.info "Starting ngircd..."

## Run ngircd
exec /etc/init.d/ngircd start